# build/ubuntu/lctrl-1.0.0/DEBIAN/postinst
#!/bin/bash

sudo systemctl daemon-reload
sudo systemctl enable lctrl
sudo systemctl restart lctrl
# sudo systemctl status lctrl

echo 'lctrl service installed and started.'